# HardenedFangbull/airootfs/etc/systemd/system/fangbull-log-optimizer.timer
[Unit]
Description=Fangbull Log Optimizer Timer
Documentation=man:fangbull-log-manager(8)
Requires=fangbull-log-optimizer.service

[Timer]
# Run every 6 hours
OnCalendar=*-*-* 00,06,12,18:00:00
# Run 5 minutes after boot
OnBootSec=5min
# Run if missed (e.g., system was off)
Persistent=true
# Add randomization to prevent all systems running at once
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
