# HardenedFangbull/airootfs/etc/systemd/system/netlink-monitor-watchdog.service
[Unit]
Description=Netlink Monitor Watchdog - Hardened Fangbull
After=network.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/bin/netlink_monitor_watchdog
Restart=on-failure
RestartSec=30
StandardOutput=null
StandardError=journal

# Improved Hardening
User=root
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=false
NoNewPrivileges=true
ReadWritePaths=/var/log/fangbull-ids
# Reduced capabilities to minimum required
CapabilityBoundingSet=CAP_KILL CAP_NET_ADMIN CAP_DAC_READ_SEARCH
# Additional security measures
RestrictNamespaces=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Resource limits
CPUQuota=10%
IOSchedulingClass=idle
MemoryLimit=100M

[Install]
WantedBy=multi-user.target 