# HardenedFangbull/airootfs/etc/systemd/system/shell-fork-bomb-terminator.service
[Unit]
Description=Shell Fork Bomb Terminator - Hardened Fangbull
After=multi-user.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/local/bin/shell_fork_bomb_terminator
Restart=on-failure
RestartSec=30
StandardOutput=null
StandardError=journal

# Improved Hardening
User=root
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ReadWritePaths=/var/log/fangbull-ids
# Reduced capabilities to minimum required
CapabilityBoundingSet=CAP_KILL CAP_DAC_READ_SEARCH
# Additional security measures
RestrictNamespaces=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Resource limits
CPUQuota=15%
IOSchedulingClass=idle
MemoryLimit=120M

[Install]
WantedBy=multi-user.target 